# scripts/products/jobs.crontab
5 */12 * * * /srv/cd-data-platform/scripts/download-products.sh
0 0 * * * /srv/cd-data-platform/scripts/cleanup-products.sh
